# user.service
[Unit]
Description = User Service
[Service]
User=roboshop
Environment=MONGO=true
<<<<<<< HEAD
Environment=REDIS_HOST=172.31.30.254
Environment=MONGO_URL="mongodb://172.31.36.162:27017/users"
=======
Environment=REDIS_HOST=172.31.84.5
Environment=MONGO_URL="mongodb://172.31.33.83:27017/users"
>>>>>>> 89cfe240350d4abeffd0c3ce9a94d2045c45db37
ExecStart=/bin/node /app/server.js
SyslogIdentifier=user

[Install]
WantedBy=multi-user.target